# root/etc/s6-overlay/s6-rc.d/pinchflat/run
#!/usr/bin/with-contenv bash

cmd="pinchflat start"
echo "*** executing => $cmd"
PHX_SERVER=true exec s6-setuidgid abc \
    $cmd
